# config/bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi -c /bbs/config/gunicorn.config.py
directory=/bbs
autostart=true
autorestart=true
