# config/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --dbpath=/data/db --port=27017 --fork
directory=/var/www/bbs
autostart=true
autorestart=true
